# packaging/systemd/utter-daemon.service
[Unit]
Description=Utter push-to-talk transcription daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/utter daemon
Restart=on-failure
RestartSec=2
# User-facing knobs (auto_paste, write_clipboard, filter_filler_words,
# PTT key) live in ~/.config/utter/config.toml. YDOTOOL_SOCKET and RUST_LOG remain here
# because they configure third-party tools, not utter itself.
Environment=YDOTOOL_SOCKET=/tmp/.ydotool_socket
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
