# packaging/systemd/utter-watcher.service
[Unit]
Description=Utter evdev key watcher (hold-to-talk)
After=graphical-session.target utter-daemon.service
Wants=utter-daemon.service
PartOf=graphical-session.target

[Service]
Type=simple
# The key comes from ~/.config/utter/config.toml (field `key`, default
# `rightmeta`). `utter set-key` updates that file. Pass `--key <name>`
# explicitly here to override the config. The 90-utter.rules udev rule
# grants read access to keyboard event devices via uaccess.
ExecStart=/usr/bin/utter watch
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
