# examples/utter-daemon.service
[Unit]
Description=Utter push-to-talk transcription daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/utter daemon
Restart=on-failure
RestartSec=2
# Always put transcribed text on the clipboard. UTTER_AUTOTYPE=1 also pastes
# it into the focused window via ydotool (uinput-based, works on any Wayland
# compositor). ydotoold's default socket path differs from its CLI's default
# when XDG_RUNTIME_DIR isn't set in the daemon env — point it explicitly at
# the system-service default path.
Environment=UTTER_AUTOTYPE=1
Environment=YDOTOOL_SOCKET=/tmp/.ydotool_socket
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
