# examples/utter-daemon.service
[Unit]
Description=Utter push-to-talk transcription daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/utter daemon
Restart=on-failure
RestartSec=2
# Write transcriptions to the primary selection; UTTER_AUTOTYPE=1 also
# synthesizes Shift+Insert via ydotool to paste into the focused window.
# UTTER_AUTOTYPE=0 skips the keystroke and leaves the user to middle-
# click or Shift+Insert manually. ydotoold's default socket path differs
# from its CLI's default when XDG_RUNTIME_DIR isn't set in the daemon
# env — point it explicitly at the system-service default path.
Environment=UTTER_AUTOTYPE=1
Environment=YDOTOOL_SOCKET=/tmp/.ydotool_socket
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
