# examples/utter-daemon.service
[Unit]
Description=Utter push-to-talk transcription daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/utter daemon
Restart=on-failure
RestartSec=2
# User-facing knobs (autotype, clipboard, cleanup, notify, PTT key) live
# in ~/.config/utter/config.toml. YDOTOOL_SOCKET and RUST_LOG remain here
# because they configure third-party tools, not utter itself. ydotoold's
# default socket path differs from its CLI's default when XDG_RUNTIME_DIR
# isn't set in the daemon env — point it explicitly at the system-service
# default path.
Environment=YDOTOOL_SOCKET=/tmp/.ydotool_socket
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
